# ansible/roles/swww/files/swww-daemon.service
[Unit]
Description=Swww Daemon
After=default.target

[Service]
Type=simple
ExecStart=/usr/local/bin/swww-daemon
Restart=always

[Install]
WantedBy=default.target
