# systemd/user/mnt-onedrive.automount
[Unit]
Description=AutoMount for /mnt/one-drive

[Automount]
Where=/mnt/onedrive
TimeoutIdleSec=600

[Install]
WantedBy=default.target
